# dist/hass-agent_0.1.0_amd64/etc/systemd/system/hass-agent.service
[Unit]
Description=HASS Agent headless service
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Default user - replace with an appropriate user for your environment
User=hassagent
WorkingDirectory=/opt/hass-agent
ExecStart=/opt/hass-agent/HASS.Agent.Headless
Restart=on-failure
RestartSec=5s
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
