# deploy/DEBIAN/prerm
#!/usr/bin/env bash
set -euo pipefail

# Stop and disable the service before removal
if command -v systemctl >/dev/null 2>&1; then
  systemctl stop hass-agent.service 2>/dev/null || true
  systemctl disable hass-agent.service 2>/dev/null || true
fi

exit 0
